# backup/deploy/backup-prune.timer
[Unit]
Description=Weekly snapshot prune

[Timer]
OnCalendar=Sun *-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
